# puppet/modules/supervisor/templates/service.ini.erb
# This file is autogenerated by Puppet. Manual changes will be overwritten!
[program:<%= name -%>]
command=<%= command %>
<% if numprocs.to_i > 1 -%>
process_name=%(program_name)s_%(process_num)02d
<% end -%>
numprocs=<%= numprocs %>
numprocs_start=<%= numprocs_start %>
priority=<%= priority %>
autostart=<%= autostart %>
autorestart=<%= autorestart %>
startsecs=<%= startsecs %>
startretries=<%= retries %>
exitcodes=<%= exitcodes %>
stopsignal=<%= stopsignal %>
stopwaitsecs=<%= stopwait %>
<% if @directory -%>
directory=<%= directory %>
<% end -%>
<% if @user -%>
user=<%= user %>
<% end -%>
redirect_stderr=<%= redirect_stderr %>
<% if @stdout_logfile -%>
stdout_logfile=<%= stdout_logfile %>
<% else -%>
stdout_logfile=/var/log/supervisor/<%= name -%>/<%= name -%>.out
<% end -%>
stdout_logfile_maxbytes=<%= stdout_logfile_maxsize %>
stdout_logfile_backups=<%= stdout_logfile_keep %>
<% if @stderr_logfile -%>
stderr_logfile=<%= stderr_logfile %>
<% else -%>
stderr_logfile=/var/log/supervisor/<%= name -%>/<%= name -%>.err
<% end -%>
stderr_logfile_maxbytes=<%= stderr_logfile_maxsize %>
stderr_logfile_backups=<%= stderr_logfile_keep %>
<% if @environment -%>
environment=<%= environment %>
<% end -%>
<% if @umask %>
umask=<%= umask %>
<% end -%>
